# compose/web/ubuntu.conf/supervisor/conf.d/php-fpm.ini
[program:php-fpm]
command=sh -c "sleep 1 && /usr/sbin/php-fpm7.4"
process_name=%(program_name)s
autostart=true
autorestart=true
startretries=5
exitcodes=0,2,70
stopsignal=QUIT
stopwaitsecs=2
stdout_logfile=/var/log/php-fpm.sup.log