# compose/web/alpine.conf/supervisor/supervisor.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.sup.log
stopasgroup=true
killasgroup=true